# ansible/.devcontainer/app/supervisord.conf
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
loglevel=info

[program:flask-app]
directory=/app
command=/opt/venv/bin/gunicorn -w 2 -b 0.0.0.0:8000 app:app
autostart=true
autorestart=true